# deploy/pidory-dev.service
[Unit]
Description=pidory-dev — Discord ↔ Claude Code Bridge (dev)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PROJECT_DIR__/target/release/pidory
Restart=on-failure
RestartSec=10

# Config — dev
Environment=PIDORY_CONFIG=__PROJECT_DIR__/config.dev.toml
Environment=HOME=__HOME_DIR__
Environment=RUST_LOG=info

EnvironmentFile=-__PROJECT_DIR__/config.dev.env
EnvironmentFile=/etc/pidory/db.env
ExecStartPre=/usr/local/bin/pidory-migrate

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=__PROJECT_DIR__ __HOME_DIR__/.claude /tmp /var/lib/pidory

[Install]
WantedBy=multi-user.target
