# deploy/pidory.service
[Unit]
Description=pidory — Discord ↔ Claude Code Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PROJECT_DIR__/target/release/pidory
Restart=on-failure
RestartSec=10

# Config
Environment=PIDORY_CONFIG=__PROJECT_DIR__/config.toml
Environment=RUST_LOG=info

EnvironmentFile=-__PROJECT_DIR__/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=__PROJECT_DIR__ __HOME_DIR__/.claude

[Install]
WantedBy=multi-user.target
